# relay-ms/conf/relay-ms.service
[Unit]
Description=relay-ms Service
After=storage.service
After=network.target

[Service]
Type=simple
EnvironmentFile=/media/fdiot/tools/env.conf
WorkingDirectory=/media/fdiot/relay-ms/bin
ExecStart=/media/fdiot/relay-ms/bin/relay-ms -ipaddr ${etcdsrv} -region ${region} -i ${devname}

#Restart=on-success
Restart=on-failure
RestartSec=6s

[Install]
WantedBy=multi-user.target